# node-server/supervisor/queue-worker.ini
[program:content-crawler]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/node app/data_source/osx-tips/queue-worker.js
numprocs=2
autostart=true
autorestart=true
user=root
pidfile=/tmp/supervisord.pid
redirect_stderr=true
directory=/var/www/content-crawler
stdout_logfile=/var/www/logs/supervisord-worker.log